# systemd/pi-control-plane.service
[Unit]
Description=Pi Control Plane Milestone 0
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi-control-plane
Group=pi-control-plane
WorkingDirectory=/var/lib/pi-control-plane
ExecStart=/usr/local/bin/pi-control-plane --config /etc/pi-control-plane/config.toml
Restart=on-failure
RestartSec=2
StateDirectory=pi-control-plane
RuntimeDirectory=pi-control-plane
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/pi-control-plane /var/log/pi-control-plane /etc/pi-control-plane

[Install]
WantedBy=multi-user.target
